# deploy/docagent-backend.service
[Unit]
Description=DocAgent Java Backend Service
After=network.target mysql.service

[Service]
Type=simple
User=docagent
Group=docagent
WorkingDirectory=/opt/docagent/backend
ExecStart=/usr/bin/java -jar /opt/docagent/backend/docagent-backend.jar --spring.profiles.active=prod
ExecStop=/bin/kill -15 $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=docagent-backend

Environment="JAVA_HOME=/usr/lib/jvm/java-21-openjdk-amd64"
Environment="DB_PASSWORD=your_strong_password"
Environment="JWT_SECRET=your_jwt_secret_key"
Environment="CORS_ORIGIN=https://your-domain.com"

[Install]
WantedBy=multi-user.target